# bot.conf
[program:my_birthday_bot]
directory=/home/ubuntu/birthdayBotNew
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/birthdayBotNew/logfile_err.log
stdout_logfile=/home/ubuntu/birthdayBotNew/logfile.log